# src/systemd/tulip-server.service
[Unit]
Description=Flask Server Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/your/repo/server
ExecStart=/path/to/uv_folder/.venv/bin/python app.py
Restart=always
RestartSec=5
Environment=FLASK_APP=app.py
Environment=FLASK_ENV=production

[Install]
WantedBy=default.target
