# tulip/systemd/tulip-check-updates.timer
[Unit]
Description=Periodically check for updates to tulip apps

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
